# patch/crontab
# OpenWrt 定时任务配置 - 2GB RAM 优化版

# 每5分钟运行系统健康检查
*/5 * * * * /usr/bin/health-monitor.sh

# 每小时清理一次日志（保留最近1000行）
0 * * * * logread | tail -n 1000 > /tmp/syslog.tmp && cat /tmp/syslog.tmp > /var/log/messages

# 每天凌晨3点重启 dnsmasq（清理DNS缓存）
0 3 * * * /etc/init.d/dnsmasq restart

# 每周日凌晨4点重启系统（可选，保持系统清爽）
# 0 4 * * 0 reboot

# 每小时检查并清理过期的连接跟踪条目
0 * * * * echo 1 > /proc/sys/net/netfilter/nf_conntrack_tcp_be_liberal 2>/dev/null

# 每30分钟同步时间（如果启用了NTP）
*/30 * * * * /usr/sbin/ntpd -q -p pool.ntp.org
